# scripts/antidistractor.service
[Unit]
Description=Antidistractor - eBPF based website blocker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/antidistractor --daemon
WorkingDirectory=/usr/local/bin
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal
# Graceful shutdown via SIGTERM → tokio ctrl_c handler
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
