# scripts/antidistractor-guard.service
[Unit]
Description=Antidistractor Guard - monitors and repairs blocking config
After=network-online.target antidistractor.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/antidistractor-guard.sh
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
